# crontab
50 8 * * * chris /srv/autosubmit/autosubmit
